# SOURCES/sshd.service
[Unit]
Description=OpenSSH server daemon (static musl build)
Documentation=man:sshd(8) man:sshd_config(5)
After=network.target sshd-keygen.service
Wants=sshd-keygen.service

[Service]
# Type=simple because this static musl build is compiled --without-systemd,
# so sshd cannot call sd_notify(). Using -D keeps sshd in the foreground.
Type=simple
EnvironmentFile=-/etc/sysconfig/sshd
ExecStartPre=/usr/bin/ssh-keygen -A
ExecStart=/usr/sbin/sshd -D $OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=42s

# Hardening
PrivateTmp=yes
ProtectHome=no
ProtectSystem=full
ReadWritePaths=/etc/ssh /var/log

[Install]
WantedBy=multi-user.target
